# infos/DESKTOP-NVM2NFS.localdomain/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (- installed on Tue Nov  4 15:25:57 2025)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)

* * * * * cd /var/www/laravel && php artisan queue:work --tries=3 --timeout=60 >> /var/www/laravel/storage/logs/queue-cron.log 2>&1
* * * * * cd /var/www/laravel && php artisan schedule:run >> /var/www/laravel/storage/logs/schedule-cron.log 2>&1
*/15 * * * * sh /var/www/laravel/bin/srv-list-nginx.sh
